# scripts/daily-stock-analysis.timer
[Unit]
Description=Daily Stock Analysis Timer

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
